# systemd/acmonitor.service
[Unit]
Description=Portable AC Monitor
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__APP_DIR__
ExecStart=__APP_DIR__/.venv/bin/acmonitor --config __APP_DIR__/config.yaml --once
